# debian/cron.d
5 5 * * * root [ -x /usr/lib/packrat/cron/cleanup.py ] && /usr/lib/packrat/cron/cleanup.py
